# config/bearings-v1.service
# systemd user-unit for the Bearings v1 rebuild.
#
# Installed (manually, by Dave) at:
#     ~/.config/systemd/user/bearings-v1.service
#
# Activated with:
#     systemctl --user daemon-reload
#     systemctl --user enable --now bearings-v1.service
#
# Concurrent run with v0.17.x:
#   * v0.17.x lives at /home/beryndil/Projects/Bearings, port 8787,
#     DB ~/.local/share/bearings/, unit `bearings.service`.
#   * v1 (this unit) lives at /home/beryndil/Projects/Bearings-v1, port 8788,
#     DB ~/.local/share/bearings-v1/, unit `bearings-v1.service`.
#
# 2026-05-04 — G1 gap-fill: ExecStart now calls `bearings serve` via the
# project venv.  The launch.py stopgap (see TODO.md "Stopgap launcher"
# entry, now resolved) has been deleted.  `bearings serve` mirrors the
# boot-time wiring launch.py performed: Settings → load_schema → create_app
# → uvicorn.run.  The DB startup hook (aiosqlite + load_schema) is in
# bearings.cli.serve._connect_db; the shutdown hook is registered via
# app.on_event("shutdown") inside _run.
#
# Historical note: at item 0.1 the FastAPI app didn't exist; this unit
# was a verifiable shape (`systemd-analyze verify` passed) until item
# 1.10.  The real app lives at `bearings.web.app:create_app` (NOT
# `bearings.api.app:app`).

[Unit]
Description=Bearings v1 (localhost web UI for Claude Code agent sessions)
Documentation=file:///home/beryndil/.claude/plans/bearings-v1-rebuild.md
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/beryndil/Projects/Bearings-v1
# Bind explicitly to localhost — Bearings is single-user / localhost only
# per the project CLAUDE.md. Uses the project venv directly (faster
# restart than `uv run`) since deps are already resolved.
ExecStart=/home/beryndil/Projects/Bearings-v1/.venv/bin/bearings serve --host 127.0.0.1 --port 8788
Restart=on-failure
RestartSec=2s

# Hardening — start permissive enough to not break uv's venv resolution but
# tight enough to fail audit-cleanly later. Tighten further in item 1.10
# when the actual API surface lands.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/bearings-v1 %h/.cache/bearings-v1 %h/.cache/uv %h/.claude %h/.local/share/claude %h/Projects
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
SystemCallArchitectures=native

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bearings-v1

[Install]
WantedBy=default.target
